# 2immerse-firmware/src/assets/DEBIAN/prerm
#!/bin/sh

# Copyright 2019 British Broadcasting Corporation
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
# 
#   http://www.apache.org/licenses/LICENSE-2.0
# 
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

set -e

if [ \( "$1" = "upgrade" -o "$1" = "remove" \) -a -L /usr/doc/parted ]; then

  # Use invoke-rc.d in dpkg maintainer script rather than calling directly
  systemctl stop kiosk || true
  systemctl stop chromium-remote-debug || true
  systemctl stop 2immerse-firmware || true
  systemctl stop dnsmasq-ap || true
  invoke-rc.d app2app-dial-server stop || true
  invoke-rc.d dvbcsstv-proxy stop || true

  pip uninstall pydvbcss -y 

  if [ -d /etc/nginx/sites-enabled -a -e /etc/nginx/sites-enabled/captive-portal.conf ]; then
    # This only removes the symbolic link
    rm -f /etc/nginx/sites-enabled/captive-portal.conf
  fi

  if [ -d /etc/nginx/sites-enabled -a ! -e /etc/nginx/sites-enabled/default -a -e /etc/nginx/sites-available/default ]; then
    # Restore default configuration
    ln -sf /etc/nginx/sites-available/default /etc/nginx/sites-enabled/default
  fi
  
  nginx -s reload || true

  # Revoke trust flags for self-signed certificate
  #certutil -d sql:$HOME/.pki/nssdb -D -n 2immerseCA

  # TODO: More clean-up is required to remove autogenerated hostapd.conf and /etc/dnsmasq.conf etc.
  #       See the configure-networking and config-ap scripts in /usr/bin
fi

